# docker/crontab
# m h dom mon dow user  command
# every hour new
30 */3 * * * /usr/bin/python3 /workspace/main.py
